# php72-strip/resources/apache.conf
[program:apache2]
command=/root/start.sh
#numprocs=1
autostart=true
autorestart=true

[program:cron]
autorestart=false
command=cron -f