# php56/resources/apache.conf
[program:apache2]
command=/root/start.sh
#numprocs=1
autostart=true
autorestart=true